# package/base-files/files/etc/uci-defaults/14_luci-settings
#!/bin/sh

uci batch <<-EOF
	set system.@system[0].zonename='Asia/Taipei'
	set system.@system[0].timezone='CST-8'
	commit system
EOF
